# ops/systemd/user/stock-insight-knowledge.service
[Unit]
Description=Stock Insight knowledge extraction and event briefs
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=1h
StartLimitBurst=3

[Service]
Type=oneshot
WorkingDirectory=/home/jigoo/.hermes/workspace/stock-insight
ExecStart=/home/jigoo/.hermes/workspace/stock-insight/apps/api/scripts/run_knowledge_pipeline.sh
TimeoutStartSec=45min
Restart=on-failure
RestartSec=10min
UMask=0077
NoNewPrivileges=true
Nice=10

[Install]
WantedBy=default.target
